# eureka-package/debian/postinst
#!/bin/bash
set -e

PKG_NAME=acu
INSTALL_DIR=/opt/acu
VENV_DIR=$INSTALL_DIR/venv
BACKEND_USER=acu
BACKEND_GROUP=acu

echo "Installing $PKG_NAME..."

# Create service user if not exists
if ! id "$BACKEND_USER" &>/dev/null; then
    useradd -r -s /bin/false -d $INSTALL_DIR -m $BACKEND_USER
    echo "Created user: $BACKEND_USER"
fi

# Setup Python venv for backend
echo "Setting up Python environment..."
mkdir -p $VENV_DIR
python3 -m venv $VENV_DIR
source $VENV_DIR/bin/activate
pip install -q -r $INSTALL_DIR/backend/requirements.txt

# Create log directory
mkdir -p /var/log/acu
chown $BACKEND_USER:$BACKEND_GROUP /var/log/acu
chmod 755 /var/log/acu

# Create runtime directory
mkdir -p /var/lib/acu
chown $BACKEND_USER:$BACKEND_GROUP /var/lib/acu
chmod 755 /var/lib/acu

# Set permissions
chown -R $BACKEND_USER:$BACKEND_GROUP $INSTALL_DIR
chmod 755 $INSTALL_DIR

# Enable systemd services
systemctl daemon-reload
systemctl enable acu-backend.service
systemctl enable acu-daemon.service

echo "$PKG_NAME installed successfully!"
echo ""
echo "Start services:"
echo "  systemctl start acu-daemon"
echo "  systemctl start acu-backend"
echo ""
echo "View logs:"
echo "  journalctl -u acu-backend -f"
echo "  journalctl -u acu-daemon -f"

exit 0
