# eureka-package/debian/prerm
#!/bin/bash
set -e

echo "Stopping ACU services..."

# Stop services
systemctl stop acu-backend.service 2>/dev/null || true
systemctl stop acu-daemon.service 2>/dev/null || true

# Disable services
systemctl disable acu-backend.service 2>/dev/null || true
systemctl disable acu-daemon.service 2>/dev/null || true

echo "Services stopped."

exit 0
